# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/opt/data/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid

; ─── Hermes Gateway (OpenAI-compatible API server on :8642) ────────────────
[program:hermes-gateway]
command=hermes gateway --platform api_server
environment=HERMES_HOME="/opt/data",API_SERVER_HOST="0.0.0.0",API_SERVER_PORT="8642",GATEWAY_ALLOW_ALL_USERS="true"
autostart=true
autorestart=true
startsecs=5
startretries=5
stdout_logfile=/opt/data/logs/gateway.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/opt/data/logs/gateway.log
priority=100

; ─── Hermes A2A Server (:9000) ────────────────────────────────────────────
[program:hermes-a2a]
command=python -m a2a_adapter
environment=HERMES_HOME="/opt/data",A2A_HOST="0.0.0.0",A2A_PORT="9000"
autostart=true
autorestart=true
startsecs=5
startretries=5
stdout_logfile=/opt/data/logs/a2a.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/opt/data/logs/a2a.log
priority=110

; ─── Filebrowser UI (:8080) ───────────────────────────────────────────────
[program:filebrowser]
command=filebrowser
    --database /opt/data/filebrowser.db
    --root /opt/data/files
    --address 0.0.0.0
    --port 8080
    --log /opt/data/logs/filebrowser.log
autostart=true
autorestart=true
startsecs=3
startretries=5
stdout_logfile=/opt/data/logs/filebrowser-supervisor.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/opt/data/logs/filebrowser-supervisor.log
priority=120
