# hermes-standard/scheduler/com.hermes.scan.plist
<?xml version="1.0"?><!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd"><plist version="1.0"><dict><key>Label</key><string>com.hermes.scan</string><key>ProgramArguments</key><array><string>/bin/sh</string><string>-c</string><string>/opt/homebrew/opt/python@3.14/bin/python3.14 /Users/rattanasak/Documents/Viber Project/Tech Tools/Hermes Agent/hermes-standard/bin/hermes_scan.py /private/tmp/claude-501/-Users-rattanasak-Documents-Viber-Project-Tech-Tools-Hermes-Agent/c1b8f77a-bebe-410d-a8eb-48f738005058/scratchpad --html /Users/rattanasak/Documents/Viber Project/Tech Tools/Hermes Agent/hermes-standard/scheduler/scan-latest.html</string></array><key>StartCalendarInterval</key><dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict></dict></plist>
