# ops/muncho/systemd/muncho-discord-connector.service.in
[Unit]
Description=Muncho privileged public Discord connector
Wants=network-online.target
After=network-online.target
Before=hermes-cloud-gateway.service

[Service]
Type=notify
NotifyAccess=main
User=muncho-discord-connector
Group=muncho-discord-connector
RuntimeDirectory=muncho-discord-connector
RuntimeDirectoryMode=0750
RuntimeDirectoryPreserve=no
StateDirectory=muncho-discord-connector
StateDirectoryMode=0700
WorkingDirectory=/
Environment=PYTHONNOUSERSITE=1
Environment=PYTHONPATH=/opt/adventico-ai-platform/hermes-agent-releases/hermes-agent-@EXACT_12_CHAR_SHA@
ExecStart=/opt/adventico-ai-platform/hermes-agent-releases/hermes-agent-@EXACT_12_CHAR_SHA@/.venv/bin/python -B -P -s -m gateway.discord_connector_bootstrap --config /etc/muncho/discord-public-connector.json
Restart=on-failure
RestartSec=5s
TimeoutStartSec=150s
TimeoutStopSec=60s
KillMode=mixed
LimitCORE=0
UMask=0077

NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
LockPersonality=yes
MemoryDenyWriteExecute=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectProc=invisible
ProtectSystem=strict
ProcSubset=pid
RemoveIPC=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
IPAddressDeny=169.254.169.254/32

ReadOnlyPaths=/opt/adventico-ai-platform/hermes-agent-releases/hermes-agent-@EXACT_12_CHAR_SHA@
ReadOnlyPaths=/etc/muncho/discord-public-connector.json
ReadOnlyPaths=/etc/muncho/discord-connector-credentials/bot-token
InaccessiblePaths=-/etc/muncho/discord-edge-credentials
InaccessiblePaths=-/var/lib/hermes-gateway/.hermes/auth.json
ReadWritePaths=/run/muncho-discord-connector
ReadWritePaths=/var/lib/muncho-discord-connector

[Install]
WantedBy=multi-user.target
